# running-systemd-inside-a-docker-container/goprinter/goprinter.service
[Unit]
Description=goprinter Service
After=syslog.target

[Service]
WorkingDirectory=/tmp
ExecStart=/appl/goprinter
#StandardOutput=journal
#StandardError=journal
